# systemd/unitfile.service.erb
[Unit]
Description=<%=name%>
Wants=multi-user.target

[Service]
ExecStart=<%=executable%>
StandardOutput=/var/log/<%=name%>.log
StandardError=/var/log/<%=name%>_error.log
Restart=always

[Install]
WantedBy=multi-user.target
